# systemd/ntfs-mount-helper.service
[Unit]
Description=NTFS Mount Helper - Corrige montagens falhas
After=local-fs.target
Before=multi-user.target
ConditionPathExists=/etc/fstab

[Service]
Type=simple
ExecStart=/usr/local/bin/ntfs-mount-helper.sh
RemainAfterExit=no
StandardOutput=journal
User=root

NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target